# quadlet/hecate-spartan.container
# Quadlet unit for hecate-spartan.
#
# Runs on Hecate INFRASTRUCTURE nodes (BEAM cluster, dedicated relay boxes,
# cooperative-contributed service nodes). NOT on user laptops.
#
# Lives in hecate-gitops/system/ and is symlinked into
# /etc/containers/systemd/ by the reconciler.

[Unit]
Description=Hecate service: hecate-spartan (mesh commons for Spartan agents)
After=network-online.target macula-station.service
Wants=network-online.target
Requires=macula-station.service

[Container]
ContainerName=hecate-spartan
Image=ghcr.io/hecate-services/hecate-spartan:latest
AutoUpdate=registry

# Realm-issued service-principal cert (read-only).
Volume=/etc/hecate/secrets/hecate-spartan:/etc/hecate/secrets:ro,Z

# Per-service state — bulk drive on the beam nodes.
Volume=/bulk0/hecate/hecate-spartan:/var/lib/hecate-spartan:Z

# Local macula-station socket (the SDK speaks to this).
Volume=/run/macula:/run/macula:Z

# Health endpoint (loopback only, for podman HEALTHCHECK).
PublishPort=127.0.0.1:8470:8470

# Entity-facing API (public — Spartan entities connect here over the network).
PublishPort=8471:8471

# Run as the dedicated `hecate` service user, never root.
User=hecate
Group=hecate

EnvironmentFile=/etc/hecate/gitops/system/hecate-spartan.env

[Service]
Restart=always
RestartSec=10s

[Install]
WantedBy=multi-user.target
